# build-config/stage-kiosk/01-kiosk-setup/files/etc/systemd/system/kiosk-display.service
[Unit]
Description=Kiosk Display Service
After=multi-user.target
Wants=disable-blanking.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/kiosk/scripts/image-viewer.sh
StandardInput=tty
TTYPath=/dev/tty1
Restart=always
RestartSec=5
Environment=TERM=linux

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
